# templates/dashboard.service.j2
[Unit]
Description=Octopus Agile Energy Dashboard
After=network.target

[Service]
Type=simple
User={{ dashboard_user }}
WorkingDirectory={{ dashboard_dir }}
ExecStart=/usr/bin/python3 -u {{ dashboard_dir }}/dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
